# roles/prometheus/templates/prometheus.service.j2
[Unit]
Description=Prometheus
After=network.target

[Service]
ExecStart=/etc/prometheus/prometheus
User=prometheus
Restart=always

[Install]
WantedBy=multi-user.target

